# cloudera-cluster-yarn-resourcemanager/supervisord.conf
[supervisord]
nodaemon=true

[program:resourcemanager]
command=yarn resourcemanager
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
