# cloudera-cluster-historyserver/supervisord.conf
[supervisord]
nodaemon=true

[program:historyserver]
command=mapred historyserver
redirect_stderr=true
stdout_logfile=/var/log/hadoop-mapreduce/mapred-historyserver.log
