# cloudera-cluster-hbase-master-backup/supervisord.conf
[supervisord]
nodaemon=true

[program:hbasemasterbackup]
command=hbase master --backup start
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0