# cloudera-cluster-zookeeper/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper]
command=/usr/lib/zookeeper/bin/zkServer.sh start-foreground
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0