# cloudera-cluster-hbase-regionserver/supervisord.conf
[supervisord]
nodaemon=true

[program:hbaseregionserver]
command=hbase regionserver start
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
